# debian/mlnx-tools.prerm
#!/bin/sh -e

deb-systemd-helper disable mlnx-eswitch.service >/dev/null  2>&1 || true

#DEBHELPER#
