# debian/mlnx-tools.prerm
systemctl disable mlnx-eswitch.service >/dev/null  2>&1 || true
