# debian/mlnx-tools.prerm
systemctl disable mlnx-bf-ctl.service >/dev/null  2>&1 || true
